# scripts/OK-CC/etc/cron.daily/data-backup-dayly
#!/bin/sh

. /etc/ipcc/cc.backup.conf

$root_path"/cc.data.backup.main.sh" >> $log_path
