# scripts/OK-CC-Clean/etc/cron.weekly/log-check-weekly
#!/bin/sh

. /etc/ipcc/cc.clean.conf

$ROOT_PATH"/cc.clean.main.sh"
